# defaults/clash/tunup.service.template
[Unit]
Description=TunUp Service
After=network.target

[Service]
ExecStart={COMMAND}
Restart=always

[Install]
WantedBy=multi-user.target
